# userpatches/xdg/autostart/astro-wizard.desktop
[Desktop Entry]
Type=Application
Name=AstroWizard
Exec=python3 /opt/astroorange/wizard/main.py --autostart
OnlyShowIn=XFCE;
