# systemd/astroorange-firstboot.service
[Unit]
Description=AstroOrange First Boot Initialization
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/astroorange-firstboot.sh
RemainAfterExit=yes
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
