# systemd/wizard-autostart.service
[Unit]
Description=Auto-start AstroOrange Wizard in VNC
After=vncserver@1.service

[Service]
Type=oneshot
User=AstroOrange
ExecStart=/usr/bin/python3 /opt/astro-wizard/app.py
Environment=DISPLAY=:1
StandardOutput=journal

[Install]
WantedBy=multi-user.target
